# systemd&supervisor/filebeat.ini
[program:filebeat]
user=huohe
command=/usr/local/software/filebeat/filebeat -c /usr/local/software/filebeat/filebeat.yml
stdout_logfile=/var/log/supervisor/filebeat.log
redirect_stderr=true
stopasgroup=true
killasgroup=true

priority=1
exitcodes=0
stopwaitsecs=3
startsecs=5
autostart = true
autorestart=true
numprocs=1
